# services/system/beo-notify-failure@.service
[Unit]
Description=Notify HA of %i failure

[Service]
Type=oneshot
ExecStart=__HOME__/beosound5c/services/system/notify-failure.sh %i
User=__USER__
