# services/system/beo-health.timer
[Unit]
Description=BeoSound 5c Health Check (every 5 min)

[Timer]
OnBootSec=120
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
