# services/system/beo-bluetooth.service
[Unit]
Description=BeoSound5c Bluetooth Remote service
After=network.target bluetooth.target
Wants=network.target bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/beosound5c/services
EnvironmentFile=-/etc/beosound5c/config.env
ExecStart=/home/pi/beosound5c/services/bluetooth.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
