# services/system/beo-health.service
[Unit]
Description=BeoSound 5c Health Check

[Service]
Type=oneshot
ExecStart=/home/kirsten/beosound5c/services/system/beo-health.sh
