# services/system/beo-masterlink.service
[Unit]
Description=BeoSound5c MasterLink Sniffer including IR data
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/beosound5c/services
EnvironmentFile=-/etc/beosound5c/config.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/beosound5c/services/masterlink.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 
