# services/system/beo-http.service
[Unit]
Description=BeoSound5c HTTP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/beosound5c/web
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m http.server 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
