# services/system/beo-spotify-fetch.service
[Unit]
Description=BeoSound5c Spotify Playlist Fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
EnvironmentFile=-/etc/beosound5c/config.env
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/pi/beosound5c/tools/spotify
ExecStart=/usr/bin/python3 /home/pi/beosound5c/tools/spotify/fetch_playlists.py
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
