# services/system/beo-ui.service
[Unit]
Description=BeoSound5c Chromium UI kiosk
After=network.target beo-http.service
Wants=network.target beo-http.service
Requires=beo-http.service

[Service]
Type=simple
User=kirsten
Group=kirsten
WorkingDirectory=/home/kirsten/beosound5c/services
ExecStart=/home/kirsten/beosound5c/services/ui.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 
